# systemd/dead-letter-replay.service
[Unit]
Description=Observability Platform - Dead Letter Replay
After=postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
User=observer
Environment="DB_DSN=dbname=observability user=observer"
ExecStart=/usr/bin/python3 /usr/local/bin/dead_letter_replay.py
StandardOutput=journal
StandardError=journal
